# deploy/systemd/cleaning-bot-backup.service
[Unit]
Description=Backup SQLite database for Cleaning Bot
Wants=network-online.target

[Service]
Type=oneshot
User=cleaningbot
Group=cleaningbot
Environment=DB_PATH=/opt/cleaning-bot/data/db.sqlite3
Environment=BACKUP_DIR=/opt/cleaning-bot/backups
Environment=RETENTION=14
ExecStart=/opt/cleaning-bot/bin/backup_sqlite.sh ${DB_PATH} ${BACKUP_DIR} ${RETENTION}
